# shell/bridge/network_start.conf
description "Initialize network bridge on startup"

start on startup
task
exec /home/bridge0/Documents/network_start.sh
