# shell/bridge/networkstart.conf
description "Initialize network bridge on startup"

start on startup
task
exec /home/bridge0/Documents/networkstart.sh
